# qr3_69.conf
[program:qr3_69]
command=python3 /home/curator3/qr3_69/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_69/test.err.log
stdout_logfile=/home/curator3/qr3_69/test.err.log
